# debian/cron.d
*/1 * * * * root [ -x /usr/lib/mcp/cron/iterate ] && /usr/lib/mcp/cron/iterate
